# gunicorn.conf
[program:gunicorn]
command=gunicorn core.wsgi:application
directory=/
user=your_username
autostart=true
autorestart=false
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log